# scripts/systemd/daphne-restart.service
[Unit]
Description=Restart Daphne Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart daphne.service
