# config/systemd/daphne.service
[Unit]
Description=Daphne ASGI server for Confio
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/confio
Environment=DJANGO_SETTINGS_MODULE=config.settings
EnvironmentFile=/opt/confio/.env
ExecStart=/opt/confio/venv/bin/daphne -b 127.0.0.1 -p 8000 config.asgi:application
Restart=always
RestartSec=3
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
