# config/systemd/telegram-ai-listener.service
[Unit]
Description=Telegram AI Listener for Confio
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=nginx
WorkingDirectory=/opt/confio
Environment=DJANGO_SETTINGS_MODULE=config.settings
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/confio/.env.mainnet
# No --chat: listen in every group the account is in. Set CONFIO_AI_TELEGRAM_ALLOWED_CHATS to restrict.
ExecStart=/opt/confio/myvenv/bin/python /opt/confio/manage.py telegram_ai_listener
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
